# init/http_exporter.service
[Unit]
Description=http_exporter
After=network.target remote-fs.target nss-lookup.target nginx.target mongod.target
[Service]
WorkingDirectory=/etc/http_exporter
ExecStart=/usr/local/bin/http_exporter --config /etc/http_exporter/config.yaml
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=http_exporter
User=root
Environment=
[Install]
WantedBy=multi-user.target

